# deploy/brightdaybot-updater.service
[Unit]
Description=BrightDayBot Safe Deploy Check
After=network-online.target brightdaybot.service
Wants=network-online.target

[Service]
Type=oneshot

# --- Customize: must match setup.sh BASE_DIR ---
Environment=BRIGHTDAYBOT_BASE=/opt/brightdaybot
ExecStart=/usr/bin/bash /opt/brightdaybot/repo/deploy/deploy.sh
# ------------------------------------------------

# Allow time for uv sync + validation
TimeoutStartSec=600
StandardOutput=journal
StandardError=journal
SyslogIdentifier=brightdaybot-deploy
