# deploy/brightdaybot.service
[Unit]
Description=BrightDayBot (Docker)
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple

# --- Customize: must match BRIGHTDAYBOT_BASE in deploy.sh ---
WorkingDirectory=/opt/brightdaybot/current
# -------------------------------------------------------------

ExecStart=/usr/bin/docker compose up --build
ExecStop=/usr/bin/docker compose down
Restart=on-failure
RestartSec=10

# Stop crash-looping: 3 failures in 5 minutes -> give up
StartLimitIntervalSec=300
StartLimitBurst=3

TimeoutStartSec=300
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
